# demo_python_service.service
[Unit]
Description=Python Demo Service

[Service]
ExecStart=/home/pi/miniconda3/envs/py36/bin/python /home/pi/dev/print_server/demo_python_service.py